# iso-profiles/antergos/live-overlay/home/antergos/.config/autostart/antergos-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Set Antergos NeXT Wallpaper
Comment=Apply Antergos NeXT default wallpaper on live session start
Exec=/usr/local/bin/set-antergos-wallpaper.sh
Terminal=false
X-KDE-autostart-phase=2
NoDisplay=true
